# fuel-bootstrap-image-builder/share/fuel-bootstrap-image/ubuntu/files.ironic/etc/init/ironic-callback.conf
description "Ironic call back script"

start on started ssh
pre-start script
  /usr/bin/configure-remote-logging.sh
end script

task

exec /usr/bin/ironic_callback
