# bootstrap/sync/etc/rc.d/rc.local
#!/bin/sh -e

mount -t devpts devpts /dev/pts || true

/usr/bin/setup-bootdev && systemctl restart network

fix-configs-on-startup || true

flock -w 1 -o /var/lock/agent.lock -c "/usr/bin/nailgun-agent >> /var/log/nailgun-agent.log 2>&1" || true

touch /var/lock/subsys/local
