# bootstrap/sync/etc/rc.d/rc.local
#!/bin/sh -e

mount -t devpts devpts /dev/pts || true

/usr/bin/setup-bootdev && systemctl restart network

fix-configs-on-startup || true

#FIXME(dteselkin): fix this in libyajl package
ln -s /usr/lib64/libyajl.so.2 /usr/lib64/libyajl.so

#/usr/bin/run-nailgun-agent
echo '* * * * * root /usr/bin/run-nailgun-agent' > /etc/cron.d/nailgun-agent

touch /var/lock/subsys/local
