# bootstrap/sync/etc/rc.d/rc.local
#!/bin/sh -e

mountpoint -q /dev/pts || mount -t devpts devpts /dev/pts

ln -nfs /usr/lib64/libyajl.so.2 /usr/lib64/libyajl.so

touch /var/lock/subsys/local
