# docker/nginx/supervisord.conf
[program:nginx]
command=/usr/sbin/nginx
autorestart=true
autostart=true
