# docker/nginx/supervisord.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autorestart=true
autostart=true
