# bootstrap/sync/etc/rc.local
#!/bin/sh -e

apt-get update

flock -w 0 -o /var/lock/agent.lock -c "/opt/nailgun/bin/agent > /var/log/agent.log 2>&1" || true
fix-mcollective-config || true

exit 0
