# iso/bootstrap_admin_node.conf
#
# This service starts bootstrap_admin_node.sh script
# which launches puppet but only if /etc/sysconfig/bootstrap_admin_node file exists.
#

start on stopped start-ttys

task
console output
script
    . /etc/sysconfig/bootstrap_admin_node
    if test "${ENABLED}" = "1"; then
        echo "ENABLED=0" > /etc/sysconfig/bootstrap_admin_node
        /usr/local/sbin/bootstrap_admin_node.sh | tee /var/log/puppet/bootstrap_admin_node.log
    fi
    initctl start tty TTY=/dev/tty1
end script
