# fuel-bootstrap-image-builder/share/fuel-bootstrap-image/ubuntu/files.ironic/etc/rc.local
#!/bin/sh -e

configure-remote-logging.sh || true

touch /var/lock/subsys/local
