# fuel-bootstrap-image-builder/share/fuel-bootstrap-image/ubuntu/files.ironic/init/ironic-callback.conf
description "Ironic call back script"

start on started ssh

task

exec /usr/local/bin/ironic_callback
