# docker/cobbler/start-container.service
[Unit]
Description=Container Startup Script
Requires=dbus.socket
After=rsyslog.service

[Service]
ExecStart=/bin/bash /usr/local/bin/start.sh > /var/log/setup.log 2>&1